# systemd/empire-backend.service
[Unit]
Description=Empire Backend API (FastAPI/Uvicorn)
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=rg
Group=rg
WorkingDirectory=/home/rg/empire-repo/backend
EnvironmentFile=/home/rg/empire-repo/backend/.env
ExecStart=/home/rg/empire-repo/backend/venv/bin/python3 -m uvicorn app.main:app --host 0.0.0.0 --port 8000 --timeout-keep-alive 65 --limit-concurrency 100
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=empire-backend

[Install]
WantedBy=multi-user.target
