# versions/v1.4.1_20260223_141332/empirebox_setup/services/empirebox-bluetooth.service
[Unit]
Description=EmpireBox Bluetooth LE Beacon
After=bluetooth.service
Wants=bluetooth.service

[Service]
Type=simple
ExecStart=/opt/empirebox/scripts/bluetooth-beacon.sh
Restart=on-failure
RestartSec=10s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=empirebox-bluetooth

[Install]
WantedBy=multi-user.target
